# .config/systemd/user/atuin-daemon.service
[Unit]
Description=Background daemon unit for Atuin

[Service]
Type=simple
ExecStart=/usr/bin/atuin daemon

[Install]
WantedBy=default.target